# liveboot_sentinel/liveboot-sentinel.service
[Unit]
Description=LiveBoot Sentinel Endpoint Agent
Documentation=https://github.com/yourorg/liveboot-sentinel
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/opt/liveboot_sentinel/agent

# Main process
ExecStart=/usr/bin/python3 /opt/liveboot_sentinel/agent/main.py

# Environment — never put secrets here; use EnvironmentFile
EnvironmentFile=/etc/liveboot_sentinel/agent.env

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/log /etc/liveboot_sentinel
ProtectHome=yes
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_DAC_READ_SEARCH

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liveboot-sentinel

[Install]
WantedBy=multi-user.target
